# collector/supervisord.conf
[program:collector]
command = perl /home/dotcloud/current/fomo_collector.pl
stderr_logfile = /var/log/supervisor/collector_error.log
stdout_logfile = /var/log/supervisor/collector.log
